# deploy/systemd/helgrind.service
[Unit]
Description=Helgrind reverse proxy gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=helgrind
Group=helgrind
WorkingDirectory=/opt/helgrind
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_ENVIRONMENT=Production
Environment=Helgrind__PublicHttpsPort=443
Environment=Helgrind__AdminHttpsPort=8444
Environment=Helgrind__DatabasePath=/var/lib/helgrind/helgrind.db
Environment=Helgrind__CertificateStoragePath=/var/lib/helgrind/certificates
EnvironmentFile=-/etc/helgrind/helgrind.env
ExecStart=/usr/bin/dotnet /opt/helgrind/Helgrind.dll
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
NoNewPrivileges=true
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
ReadOnlyPaths=/opt/helgrind
StateDirectory=helgrind
StateDirectoryMode=0750
ReadWritePaths=/var/lib/helgrind

[Install]
WantedBy=multi-user.target